# dotfiles_gitea/config/goxlr/GoXLR.service
[Unit]
Description=GoXLR Utility Daemon
After=network.target

[Service]
ExecStart=/usr/bin/goxlr-daemon

[Install]
WantedBy=default.target